# etc/cron.d/maas-import-pxe-files
# Runs the MAAS image downloads early every Sunday.

11 4 * * 0 root	/usr/sbin/maas-import-pxe-files &> /dev/null
