# etc/cron.d/maas-import-isos
# Runs maas-import-isos every Sunday at 4am
# This will get the new mini.iso and generate a new cobbler profile
11 4 * * 0 root  /usr/sbin/maas-import-isos &> /dev/null
