# etc/cron.d/maas-import-pxe-files
# Runs the MAAS image downloads early every Sunday.

# Installation images:
11 4 * * 0 root	/usr/sbin/maas-import-pxe-files &> /dev/null

# Ephemeral images:
44 4 * * 0 root /usr/sbin/maas-import-ephemerals &> /dev/null
