# debian/maas-region-controller.prerm
#!/bin/sh

set -e

. /usr/share/debconf/confmodule

#DEBHELPER#

invoke-rc.d maas-regiond stop || true

. /usr/share/dbconfig-common/dpkg/prerm.pgsql
dbc_go maas-region-controller $@
