# debian/maas-rack-controller.postinst
#!/bin/sh

set -e

. /usr/share/debconf/confmodule
db_version 2.0

configure_logging() {
    # Give appropriate permissions
    if [ ! -f /var/log/maas/rackd.log ]; then
        touch /var/log/maas/rackd.log
    fi
    chown maas:maas /var/log/maas/rackd.log
}

configure_libdir() {
    # If /var/lib/maas/secret exists, ensure ownership and
    # permissions are correct
    if [ -f /var/lib/maas/secret ]; then
        chown maas:maas /var/lib/maas/secret
        chmod 0640 /var/lib/maas/secret
    fi

    if [ -f /var/lib/maas/maas_id ]; then
       chown maas:maas /var/lib/maas/maas_id
    fi

    # Ensure dhcpd related files have the correct permissions
    if [ -f /var/lib/maas/dhcpd.conf ]; then
        chown root:maas /var/lib/maas/dhcpd.conf
    fi
    if [ -f /var/lib/maas/dhcpd-interfaces ]; then
        chown root:maas /var/lib/maas/dhcpd-interfaces
    fi
    if [ -f /var/lib/maas/dhcpd6.conf ]; then
        chown maas:maas /var/lib/maas/dhcpd6.conf
    fi
    if [ -f /var/lib/maas/dhcpd6-interfaces ]; then
        chown maas:maas /var/lib/maas/dhcpd6-interfaces
    fi
    if [ -S /var/lib/maas/dhcpd.sock ]; then
        chown maas:maas /var/lib/maas/dhcpd.sock
    fi
}

extract_cluster_uuid(){
    # Extract ClUSTER_UUID setting from config file $1. This will work
    # on the old the cluster config file (which is shell).
    awk '{ split($1,array,"\"")} END{print array[2] }' "$1"
}

extract_maas_url(){
    # Extract the MAAS_URL setting from the config file $1.
    grep -Eo "https?://[^ ]+" "$1" | cut -d"\"" -f1
}

configure_maas_url(){
    # Get the MAAS URL on configure/reconfigure and write it to the conf files.
    db_get maas-rack-controller/maas-url || true
    if [ -n "$RET" ]; then
        maas-rack config --region-url "$RET"
    fi
}

configure_shared_secret() {
    db_get maas-rack-controller/shared-secret || true
    if [ -n "$RET" ]; then
        echo "$RET" | maas-rack install-shared-secret
        chown maas:maas /var/lib/maas/secret
        chmod 0640 /var/lib/maas/secret
    fi
}

edit_named_options() {
    # Remove any existing MAAS-related include line from
    # /etc/bind/named.conf.local, then re-add it.
    sed -i '/^include\s.*maas/d' /etc/bind/named.conf.local
    /usr/lib/maas/maas-common get-named-conf --edit --config-path /etc/bind/named.conf.local

    # Add a line in /etc/bind/named.conf.options that includes the
    # /etc/named/maas/named.conf.options.inside.maas file.
    /usr/lib/maas/maas-common edit-named-options --config-path /etc/bind/named.conf.options
}

fix_dns_permissions() {
    if [ -d /etc/bind/maas ]; then
        chown maas:root /etc/bind/maas
        chown -R maas:maas /etc/bind/maas/*
    fi
    if [ -f /etc/bind/maas/named.conf.maas ]; then
        chown maas:maas /etc/bind/maas/named.conf.maas
        chmod 644 /etc/bind/maas/named.conf.maas
    fi
    if [ -f /etc/bind/maas/named.conf.options.inside.maas ]; then
        chown maas:maas /etc/bind/maas/named.conf.options.inside.maas
        chmod 644 /etc/bind/maas/named.conf.options.inside.maas
    fi
    if [ -f /etc/bind/maas/rndc.conf.maas ]; then
        chown maas:root /etc/bind/maas/rndc.conf.maas
        chmod 600 /etc/bind/maas/rndc.conf.maas
    fi
    if [ -f /etc/bind/maas/named.conf.rndc.maas ]; then
        chown maas:bind /etc/bind/maas/named.conf.rndc.maas
        chmod 640 /etc/bind/maas/named.conf.rndc.maas
    fi
}

# Unconditionally ensure that there is at least an empty configuration
# file. This does *not* overwrite any existing configuration.
maas-rack config
chown maas:maas /etc/maas/rackd.conf

if [ "$1" = "configure" ] && [ -z "$2" ]; then
    configure_logging
    configure_libdir
    configure_maas_url
    configure_shared_secret
    # Configure DNS
    # If /etc/bind/maas is empty, setup dns.
    if [ ! "$(ls -A /etc/bind/maas)" ]; then
        /usr/lib/maas/maas-common setup-dns
    fi
    fix_dns_permissions
    edit_named_options

elif [ -n "$DEBCONF_RECONFIGURE" ]; then
    configure_maas_url
    # Only ask for a shared secret when the region is not installed
    # on the same system.
    if [ ! -f /usr/sbin/maas-region ]; then
        db_input high maas-rack-controller/shared-secret
        db_go
    fi
    configure_shared_secret

elif [ "$1" = "configure" ] && dpkg --compare-versions "$2" gt 0.1+bzr266+dfsg-0ubuntu1; then
    configure_logging
    configure_libdir
    # Configure DNS
    # If /etc/bind/maas is empty, setup dns.
    if [ ! "$(ls -A /etc/bind/maas)" ]; then
        /usr/lib/maas/maas-common setup-dns
    fi
    fix_dns_permissions
    edit_named_options
fi

deb-systemd-helper enable 'maas-rackd.service' >/dev/null || true
deb-systemd-invoke restart 'maas-rackd.service' >/dev/null || true
deb-systemd-invoke restart 'bind9.service' >/dev/null || true

db_stop

#DEBHELPER#
