# debian/maas-region-api.prerm
#!/bin/sh

set -e

. /usr/share/debconf/confmodule

#DEBHELPER#

invoke-rc.d apache2 stop || true
